# scripts/com.moneyflow.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  자금 조류 데일리 launchd 에이전트 (참조용 템플릿 · 비밀키 없음)
  ---------------------------------------------------------
  설치:  cp scripts/com.moneyflow.daily.plist ~/Library/LaunchAgents/
         launchctl load ~/Library/LaunchAgents/com.moneyflow.daily.plist
  해제:  launchctl unload ~/Library/LaunchAgents/com.moneyflow.daily.plist

  평일(월~금) 08:00 (시스템 로컬시간 = KST) 1회 실행.
  비밀키는 plist가 아니라 scripts/run_daily.env에서 로드된다.
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.moneyflow.daily</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>/Users/dkmac/Desktop/@26/money_flow_2606/scripts/run_daily.sh</string>
    </array>

    <!-- 평일 08:00 (월=1 … 금=5) -->
    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/dkmac/Desktop/@26/money_flow_2606/scripts/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dkmac/Desktop/@26/money_flow_2606/scripts/launchd.err.log</string>
  </dict>
</plist>
